# docker/web/supervisor/php-fpm.conf
[program:php-fpm]
command=php-fpm7.1
numprocs=1
autostart=true
autorestart=true
redirect_stderr=true
